# agents/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
childlogdir=/var/log/supervisor

# Docker daemon - must start first (priority 1) for Docker-in-Docker
[program:dockerd]
command=dockerd --storage-driver=overlay2
priority=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dockerd.log
stderr_logfile=/var/log/supervisor/dockerd_err.log
startsecs=5

# PostgreSQL database - start early (priority 2)
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
user=postgres
priority=2
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log
startsecs=5

# Nginx reverse proxy and static file server (priority 50)
[program:nginx]
command=nginx -g "daemon off;"
priority=50
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
startsecs=2

[program:software_engineering_api]
command=python -m uvicorn software_engineering_team.api.main:app --host 0.0.0.0 --port 8000
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sw_api.log
stderr_logfile=/var/log/supervisor/sw_api_err.log
environment=PYTHONPATH="/app",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"

[program:blogging_api]
command=python -m uvicorn blogging.api.main:app --host 0.0.0.0 --port 8001
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/blogging_api.log
stderr_logfile=/var/log/supervisor/blogging_api_err.log
environment=PYTHONPATH="/app:/app/blogging",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"

[program:market_research_api]
command=python -m uvicorn market_research_team.api.main:app --host 0.0.0.0 --port 8002
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/market_research_api.log
stderr_logfile=/var/log/supervisor/market_research_api_err.log
environment=PYTHONPATH="/app",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"

[program:soc2_compliance_api]
command=python -m uvicorn soc2_compliance_team.api.main:app --host 0.0.0.0 --port 8003
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/soc2_compliance_api.log
stderr_logfile=/var/log/supervisor/soc2_compliance_api_err.log
environment=PYTHONPATH="/app",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"

[program:social_marketing_api]
command=python -m uvicorn social_media_marketing_team.api.main:app --host 0.0.0.0 --port 8004
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/social_marketing_api.log
stderr_logfile=/var/log/supervisor/social_marketing_api_err.log
environment=PYTHONPATH="/app",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"

[program:blog_research_api]
command=python -m uvicorn api.main:app --host 0.0.0.0 --port 8005
directory=/app
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/blog_research_api.log
stderr_logfile=/var/log/supervisor/blog_research_api_err.log
environment=PYTHONPATH="/app:/app/blogging",NVM_DIR="/usr/local/nvm",PATH="/usr/local/nvm/versions/node/v22.12/bin:%(ENV_PATH)s"
